# systemd.service
[Unit]
Description=Unknown CRUD Backend API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/api
Environment=PATH=/var/www/html/api/env/bin
Environment=PYTHONPATH=/var/www/html/api
ExecStart=/var/www/html/api/env/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target